# periodic/daily/450.letsencrypt-certbot
/usr/local/bin/certbot renew --quiet
/usr/local/etc/rc.d/nginx restart
